# infra-ansible/roles/ansible_pull/templates/ansible-pull.timer.j2
[Unit]
Description=Run ansible-pull periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ ansible_pull_interval }}
Persistent=true

[Install]
WantedBy=timers.target